# src/deb/prerm
#!/bin/sh

if [ "$1" = "remove" ]; then
    invoke-rc.d teamcity-agent stop || true
fi

exit 0

